# eos-remove-old-flatpak-extension-dir.service
[Unit]
Description=Remove old Flatpak unmanaged extensions directory from /var/lib/flatpak
DefaultDependencies=false
Conflicts=shutdown.target
Before=systemd-tmpfiles-setup.service
After=local-fs.target
Wants=local-fs.target

# Only run on ARM systems where the old directory exists
ConditionArchitecture=arm
ConditionPathIsDirectory=/var/lib/flatpak/extension
ConditionPathIsSymbolicLink=!/var/lib/flatpak/extension

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/rm -rf /var/lib/flatpak/extension

[Install]
WantedBy=multi-user.target
