# flatpak-extension-runtime-fdo-openh264-watcher.path.in
[Unit]
After=systemd-tmpfiles-setup.service

[Path]
PathChanged=/var/lib/flatpak/runtime/org.freedesktop.Platform.openh264/@FLATPAK_ARCH@/2.2.0/active/files/extra/libopenh264.so.6

[Install]
WantedBy=multi-user.target
