# dracut/repartition/endless-repartition.service
[Unit]
Description=EOS repartitioning

# Dependencies based around dracut-pre-mount.service
DefaultDependencies=no
Before=initrd-root-fs.target sysroot.mount
After=dracut-pre-mount.service
ConditionPathExists=/etc/initrd-release

# Avoid running at the same time as fsck, and also use that unit's completion
# as a synchronization point to know that the root device is ready.
After=systemd-fsck-root.service

[Service]
Type=oneshot
ExecStart=-/bin/endless-repartition
RemainAfterExit=yes
