# factory-test/eos-factory-test.target
[Unit]
Description=Endless factory test target
AllowIsolate=yes
Requires=multi-user.target
After=multi-user.target
